# configuration/init/openrc/c8y-firmware-plugin
#!/sbin/openrc-run
description="Thin-edge device firmware management for Cumulocity"
command="/usr/bin/c8y-firmware-plugin"
command_background="yes"
pidfile="/run/lock/${RC_SVCNAME}.lock"
command_user="tedge:tedge"

depends() {
    after mosquitto
}
