# configuration/init/initd/c8y-configuration-plugin
#!/sbin/openrc-run
description="Thin-edge device configuration management for Cumulocity"
command="/usr/bin/c8y-configuration-plugin"
command_background="yes"
pidfile="/run/lock/${RC_SVCNAME}.lock"

depends() {
    after mosquitto
}
